# image/patches/root-fs/common/etc/systemd/system/serial-getty@ttyGS0.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin tf --keep-baud 115200,38400,9600 %I xterm-color
