# image/patches/root-fs/common/lib/udev/rules.d/40-rtl8821cu.rules
echo ATTR{idVendor}=="0bda", ATTR{idProduct}=="1a2b", RUN+="/usr/sbin/usb_modeswitch -KW -v 0bda -p 1a2b"
