# image/patches/root-fs/common/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin tf --noclear %I $TERM
